# infra/panchanga-backend.service
[Unit]
Description=Panchanga Backend API
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/apps/panchanga/backend
Environment=PATH=/apps/panchanga/backend/venv/bin:/usr/bin:/bin
ExecStart=/apps/panchanga/backend/venv/bin/uvicorn api:app --host 127.0.0.1 --port 8121
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
